# recipes-raspi-zero-cfg/systemd/files/wifi.service
[Unit]
Description=wifi

[Service]
ExecStart=modprobe brcmfmac

[Install]
WantedBy=multi-user.target
